# dutagent-configs/dutagent.service
[Unit]
Description=DUT Agent Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/tmp
ExecStart=/usr/local/bin/dutagent -a 0.0.0.0:1024 -c /etc/dutagent/config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target